# qbittorrent/rootfs/etc/services.d/qbittorrent/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "Starting qBittorrent..."

# в”Ђв”Ђв”Ђ Read options from HA Configuration panel в”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђ
WEBUI_PORT=$(bashio::config 'WEBUI_PORT')
TORRENTING_PORT=$(bashio::config 'TORRENTING_PORT')
PUID=$(bashio::config 'PUID')
PGID=$(bashio::config 'PGID')
TZ=$(bashio::config 'TZ')

# в”Ђв”Ђв”Ђ Persistent config (survives add-on updates) в”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђ
CONFIG_DIR="/share/qbittorrent"
QB_CONFIG="${CONFIG_DIR}/config"
mkdir -p "${QB_CONFIG}"
mkdir -p "${CONFIG_DIR}/webgui"

bashio::log.info "Config dir : ${CONFIG_DIR}"
bashio::log.info "WebUI port : ${WEBUI_PORT}"
bashio::log.info "Torrent port: ${TORRENTING_PORT}"

# в”Ђв”Ђв”Ђ Create user-configured media folders в”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђ
FIRST_FOLDER=""
for folder in $(bashio::config 'media_folders'); do
    if [[ "${folder}" != /media/* ]] && [[ "${folder}" != /share/* ]]; then
        bashio::log.warning "Skipping '${folder}': only /media/* and /share/* paths are allowed."
        continue
    fi
    mkdir -p "${folder}"
    bashio::log.info "Folder ready: ${folder}"
    [[ -z "${FIRST_FOLDER}" ]] && FIRST_FOLDER="${folder}"
done

if [[ -z "${FIRST_FOLDER}" ]]; then
    bashio::log.warning "No valid folders configured вЂ” falling back to /media/downloads"
    FIRST_FOLDER="/media/downloads"
    mkdir -p "${FIRST_FOLDER}"
fi
bashio::log.info "Default save path: ${FIRST_FOLDER}"

# в”Ђв”Ђв”Ђ Bootstrap qBittorrent.conf в”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђ
QB_CONF="${QB_CONFIG}/qBittorrent/qBittorrent.conf"
if [ ! -f "${QB_CONF}" ]; then
    bashio::log.info "First run вЂ” creating qBittorrent.conf"
    mkdir -p "$(dirname "${QB_CONF}")"
    cat > "${QB_CONF}" << EOF
[AutoRun]
enabled=false

[BitTorrent]
Session\DefaultSavePath=${FIRST_FOLDER}
Session\TempPath=${FIRST_FOLDER}/temp
Session\TempPathEnabled=false
Session\Port=${TORRENTING_PORT}
Session\Interface=
Session\InterfaceName=

[Core]
AutoDeleteAddedTorrentFile=Never

[LegalNotice]
Accepted=true

[Meta]
MigrationVersion=6

[Preferences]
Connection\PortRangeMin=${TORRENTING_PORT}
Downloads\SavePath=${FIRST_FOLDER}
Downloads\TempPath=${FIRST_FOLDER}/temp
Downloads\TempPathEnabled=false
General\Locale=en
General\UseRandomPort=false
WebUI\Address=*
WebUI\AlternativeUIEnabled=false
WebUI\BanDuration=3600
WebUI\CSRFProtection=true
WebUI\ClickjackingProtection=true
WebUI\HTTPS\Enabled=false
WebUI\HostHeaderValidation=true
WebUI\LocalHostAuth=true
WebUI\MaxAuthenticationFailCount=5
WebUI\Port=${WEBUI_PORT}
WebUI\RootFolder=
WebUI\SecureCookie=true
WebUI\ServerDomains=*
WebUI\SessionTimeout=3600
WebUI\Username=admin
EOF
    bashio::log.info "Config created. Login: admin / adminadmin  в†ђ CHANGE THE PASSWORD!"
else
    sed -i "s|^Downloads\\\\SavePath=.*|Downloads\\\\SavePath=${FIRST_FOLDER}|" "${QB_CONF}" || true
    sed -i "s|^Session\\\\DefaultSavePath=.*|Session\\\\DefaultSavePath=${FIRST_FOLDER}|" "${QB_CONF}" || true
    bashio::log.info "Config updated."
fi

# в”Ђв”Ђв”Ђ Run qBittorrent as the specified user в”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђв”Ђ
bashio::log.info "Launching qbittorrent-nox (PUID=${PUID} PGID=${PGID})..."

# Ensure user/group exist
if ! getent group "${PGID}" > /dev/null 2>&1; then
    addgroup -g "${PGID}" qbittorrent 2>/dev/null || true
fi
if ! getent passwd "${PUID}" > /dev/null 2>&1; then
    adduser -D -H -u "${PUID}" -G qbittorrent qbittorrent 2>/dev/null || true
fi

chown -R "${PUID}:${PGID}" "${CONFIG_DIR}" 2>/dev/null || true

exec su-exec "${PUID}:${PGID}" \
    qbittorrent-nox \
    --profile="${QB_CONFIG}" \
    --webui-port="${WEBUI_PORT}"
